# deploy/immersion-scan.service
[Unit]
Description=Immersion: queue a library scan

[Service]
Type=oneshot
WorkingDirectory=%h/dev/immersion/server
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin
ExecStart=%h/.local/bin/uv run immersion scan
