# deploy/immersion-api.service
[Unit]
Description=Immersion API (uvicorn on 127.0.0.1:8770)
After=network.target

[Service]
WorkingDirectory=%h/dev/immersion/server
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin
ExecStart=%h/.local/bin/uv run immersion api
Restart=on-failure

[Install]
WantedBy=default.target
