# deploy/immersion-scan.timer
[Unit]
Description=Immersion: periodic library scan

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
